# low_prio_machine.slice
[Unit]
Description=Low Priority Virtual Machine and Container Slice
Documentation=man:systemd.special(7)
Before=slices.target
[Slice]
CPUSetMems=all
CPUSetCpus=all
CPUSetCloneChildren=1
MemoryAccounting=true
CPUAccounting=true
FreezerAccounting=true
BlockIOAccounting=true